# deploy/verdict-backend.service
[Unit]
Description=Verdict FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/verdict-backend
EnvironmentFile=/opt/verdict-backend/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/verdict-backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001 --workers 2 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=5
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/verdict-backend

[Install]
WantedBy=multi-user.target
